# deploy/supervisord-xhs.conf
[program:xhs-mcp]
command=/opt/xhs-mcp/start.sh
directory=/opt/xhs-mcp
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/xhs-mcp.log
stderr_logfile=/var/log/xhs-mcp-error.log
environment=DISPLAY=":1"
